# docker/supervisord.conf
; supervisord configuration for SteamPulse Docker container
; Manages three processes: nginx (web server), scheduler (periodic fetch loop),
; and steam-serve (mutation API + re-render sidecar on 127.0.0.1:8081)

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
priority=10

[program:scheduler]
command=/scheduler.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
priority=20

[program:steam-serve]
command=steam-serve --config /run/steampulse/config.toml --db /data/steam_library.db --output-dir /data --host 127.0.0.1 --port 8081
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
priority=30
